# scripts/connman
#!/sbin/runscript
##
## Author(s):
##  - Samuel MARTIN <smartin@aldebaran-robotics.com>
##
## Copyright (C) 2010 - 2011 Aldebaran Robotics
##

BINARY="/etc/init.d/connman.sh"

depend() {
  need dbus
}

start() {
  if [ -f /home/nao/naoqi/config/connman.profile ]; then
  	echo "Applying profile /home/nao/naoqi/config/connman.profile"
        rm /var/lib/connman/default.profile
  	cp /home/nao/naoqi/config/connman.profile /var/lib/connman/default.profile
  fi
  "${BINARY}" &
}

stop() {
  killall -9 connman.sh
  killall -9 connmand
}

restart() {
  stop
  sleep 2
  start
}

